# supervisord.conf

[program:nwc2]
command=python2 /nwc2/main.py 889%(process_num)01d
process_name=%(program_name)s_%(process_num)01d
redirect_stderr=true
stdout_logfile=/var/log/nwc2.log
numprocs=4
